# deploy/scripts/supervisor.conf
# ================================================================
# Supervisor 配置 — WebChat 后端进程管理
# 适用于无法使用 Systemd 的环境（如 Docker 容器内）
# 安装: sudo apt install supervisor
# 用法: supervisorctl start webchat
# ================================================================

[program:webchat]
command=java -Xms256m -Xmx512m -XX:+UseG1GC -jar /opt/webchat/app.jar
directory=/opt/webchat
user=webchat
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/webchat/stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/webchat/stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=
    MYSQL_HOST="127.0.0.1",
    MYSQL_PORT="3306",
    MYSQL_USER="root",
    MYSQL_PASSWORD="webchat2026",
    MYSQL_DATABASE="webchat",
    MONGO_HOST="127.0.0.1",
    MONGO_PORT="27017",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    JWT_SECRET="change-this-to-strong-secret"
